# vultr-whisper.timer
[Unit]
Description=Run Vultr Whisper transcription daily at 5 PM

[Timer]
OnCalendar=Mon..Thu 18:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
